# menu@.service
.include /lib/systemd/system/getty@.service

[Unit]
Description=htop on %I
Documentation=man:htop(1)

[Service]
Environment=
Environment=TERM=linux HOME=/root
ExecStart=
ExecStart=/home/victor/111
StandardInput=tty-fail
StandardOutput=tty

